# overlay/etc/crontabs/root
# min   hour    day     month   weekday command
*       *       *       *       *       bf-cron -q /etc/periodic/1min
*/15    *       *       *       *       bf-cron /etc/periodic/15min
0       *       *       *       *       bf-cron /etc/periodic/hourly
0       2       *       *       *       bf-cron /etc/periodic/daily
0       3       *       *       6       bf-cron /etc/periodic/weekly
0       5       1       *       *       bf-cron /etc/periodic/monthly
