# overlay/etc/services.d/cron-log/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Forward Cron error log to Docker
#======================================================================================================================

FILE=/var/log/cron.log && touch ${FILE}
bcg-forward "cron-log" ${FILE}
