# overlay/etc/s6-overlay/s6-rc.d/00-init/up
foreground
{
    bf-fix-attrs
    bf-cont-init
}