# overlay/etc/services.d/cron/run
#!/usr/bin/with-contenv sh
_echo "Starting cron"
/usr/sbin/crond -f
