# overlay/etc/services.d/cron/run
#!/bin/sh

_e " - starting cron service"
/usr/sbin/crond -f
