# overlay/etc/services.d/cron/run
#!/bin/sh
echo $PATH
_echo " - starting cron service"
/usr/sbin/crond -f
