# overlay/etc/services.d/cron/run
#!/usr/bin/with-contenv sh
_echo "starting cron service"
/usr/sbin/crond -f
