# overlay/etc/services.d/cron/run
#!/bin/sh
_echo " - starting cron service"
/usr/sbin/crond -f
