# overlay/etc/services.d/00-crond/run
#!/bin/bash

echo "Registering cron daemon..."
rc-service crond start && rc-update add crond
