# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    ifelse { bf-fix-attrs } { bf-cont-init } bf-svc-terminate
}
