# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { /usr/lib/bf/init }
}
