# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    /usr/local/bin/bf-fix-attrs
    /usr/local/bin/bf-cont-init
}
