# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { bf-echo "Initialising container." "init" }
    if { /usr/lib/bf/init || bf-cont-terminate }
    if { bf-done "init" }
}
