# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    ifelse -n { /command/with-contenv bf-init }
    {
        if { bf-cont-terminate init/up }
    }
}
