# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { bf-fix-attrs }
    if { bf-cont-init }
}
