# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { ${BF_LIB}/init }
}
