# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { bf-echo "Initialising container." "init" }
    ifelse { /usr/lib/bf/init } { bf-error }
    if { bf-done "init" }
}
