# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { chmod +x /usr/local/bin/bf-fix-attrs }
    if { bf-fix-attrs }
    if { bf-cont-init }
}
