# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { bf-echo "Initialising container." "init" }
    ifelse -n { /usr/lib/bf/init } { bf-notok "Unable to initialise container." && bf-cont-terminate }
    if { bf-done "init" }
}
