# overlay/etc/s6-overlay/s6-rc.d/s6-init/up
foreground
{
    if { bf-s6-init-env }
}
