# overlay/etc/s6-overlay/s6-rc.d/s6-init/up
foreground
{
    if { bf-debug "S6 Init" }
}
