# overlay/etc/s6-overlay/s6-rc.d/s6-init/up
foreground
{
    if { echo "S6_CONT_ENV=$S6_CONT_ENV" }
    if { ls $S6_CONT_ENV }
}
