# overlay/etc/s6-overlay/s6-rc.d/cron/finish
#!/command/with-contenv nu

use bf
use s6
s6 env set_executable cron

# Show helpful log message
def main [...args] {
    s6 svc finish --terminate
    bf clean
}
