# overlay/etc/s6-overlay/s6-rc.d/cron/finish
#!/usr/bin/nu

use bf
use bf-s6
bf env load -e cron

# Show helpful log message
def main [...args] { bf-s6 svc finish --terminate }
