# overlay/etc/s6-overlay/s6-rc.d/cron/run
#!/usr/bin/nu-env

use bf
use s6
s6 env set_executable cron

# Start cron daemon
def main [...args] {
    # crond flags
    #   -f  run in foreground mode (allows S6 to supervise the service - which keeps the container running)
    #   -l  log level (from 0 to 8 - 0 is the most verbose)
    #   -L  log to a file (instead of stdout)
    bf write "Starting cron."
    /usr/sbin/crond -f -l (bf env req CRON_LOG_LEVEL) -L (bf env req CRON_LOG_FILE)
}
