# overlay/etc/profile.d/path.sh
#!/bin/sh

export PATH=/command:${PATH}
